# systemd/npm-template-sync-github-hook.socket
[Socket]
ListenStream=3100

[Install]
WantedBy=sockets.target
