# systemd/npm-template-sync-github-hook.socket
[Socket]
ListenStream=3100
MaxConnections=4

[Install]
WantedBy=sockets.target
