# systemd/npm-template-sync-github-hook.socket
[Socket]
ListenStream=3100
BindIPv6Only=both
MaxConnections=4

[Install]
RequiredBy={{name}}.service
