# systemd/npm-template-sync-github-hook.service
[Unit]
Description={{description}}
After=network-online.target

[Service]
Type=notify

ExecStart={{installdir}}/bin/{{name}}
TimeoutStartSec=20
StandardOutput=syslog
DynamicUser=yes
NoNewPrivileges=true
ConfigurationDirectory={{name}}
RuntimeDirectory={{name}}

[Install]
WantedBy=multi-user.target
