# systemd/npm-template-sync-github-hook.service
[Unit]
Description={{description}}

[Service]
Type=notify
#Environment="GH_TOKEN=XXXX"
#Environment="WEBHOOK_SECRET=YYYY"

ExecStart={{installdir}}/bin/{{name}}
TimeoutStartSec=20
StandardOutput=syslog
User=nobody
Group=nobody

RuntimeDirectory={{name}}

[Install]
WantedBy=multi-user.target
