# systemd/npm-template-sync-github-hook.service
[Unit]
Description={{description}}
After=network-online.target

[Service]
Type=notify

ExecStart={{installdir}}/bin/{{name}}
TimeoutStartSec=20
StandardOutput=syslog
User=nobody
Group=nobody
PrivateTmp=true
NoNewPrivileges=true
RuntimeDirectory={{name}}

[Install]
WantedBy=multi-user.target
