# systemd/npm-template-sync-github-hook.socket
[Socket]
ListenStream=/run/{{name}}/socket
MaxConnections=2

[Install]
RequiredBy={{name}}.service
