# systemd/npm-template-sync-github-hook.socket
[Socket]
ListenStream={{socket}}
MaxConnections=2

[Install]
RequiredBy={{name}}.service
